# init-files/supervisord_sendmail.ini
[program:sendmail]
command=/usr/sbin/sendmail -bd
autostart=true
autorestart=false
numprocs=1
exitcodes=0
startretries=0
